# deploy/pisonet.service
[Unit]
Description=PisoNet Server
Documentation=https://github.com/yourrepo/pisonet
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pisonet/server
Environment="PATH=/home/pi/pisonet/server/venv/bin"
ExecStart=/home/pi/pisonet/server/venv/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-level info

# Always restart on crash (3 second cooldown)
Restart=always
RestartSec=3

# Kill cleanly
KillMode=mixed
TimeoutStopSec=10

# Log to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pisonet

[Install]
WantedBy=multi-user.target
